# templates/jellyfin.service.epp
<%- | Optional[String]  $system_user,
| -%>
[Unit]
Description=Jellyfin
After=network.target

[Service]
Type=simple
User=<%= $system_user %>
Restart=always
ExecStart=/opt/jellyfin/start.sh

[Install]
WantedBy=multi-user.target
